# etc/docker/supervisor/apache2.conf
[program:apache2]
command=/bin/bash -c "service php8.1-fpm start && /usr/sbin/apache2ctl -D FOREGROUND"